# supervisord.conf
[supervisord]
nodaemon=true

[program:hdfs]
command=/bin/bash -c "bootstrap.sh"